# etc/cron.d/slack_zoom_check
export TERM=xterm-termite
* * * * * ghislain /usr/local/bin/update_zoom_slack_status
